# 2floppy/local/inittab

id:1:initdefault:  

# System initialization.
si::sysinit:/etc/rc.d/rc.sysinit

1:12345:respawn:/bin/bash
